# root/etc/s6-overlay/s6-rc.d/saslauthd/run
#!/usr/bin/with-contenv bash

case ${USE_SASLAUTHD} in
    yes|true)
	exec /usr/sbin/saslauthd -a sasldb -d
	;;
    *)
	;;
esac


