# root/etc/services.d/opendkim/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/opendkim -f