# root/etc/s6-overlay/s6-rc.d/opendkim/run
#!/usr/bin/with-contenv bash

case ${USE_DKIM} in
    yes|true)
	exec /usr/sbin/opendkim -f
	;;
    *)
	;;
esac

